# vagrant/chef/cookbooks/mysql/templates/default/omnios/mysql.xml.erb
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='database/mysql' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <exec_method name='start' type='method' exec='/lib/svc/method/mysqld start' timeout_seconds='18446744073709551615'/>
    <exec_method name='stop' type='method' exec='/lib/svc/method/mysqld stop' timeout_seconds='18446744073709551615'/>
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>MySQL RDBMS</loctext>
      </common_name>
      <documentation>
        <manpage title='MySQL <%= @version %>' section='1'/>
        <doc_link name='mysql.com' uri='http://dev.mysql.com/docs'/>
      </documentation>
    </template>
  </service>
</service_bundle>
